# microwhoisd-init
#!/sbin/openrc-run

command="${MICROWHOISD_BIN}"
command_args="${MICROWHOISD_OPTS}"
command_background=yes
pidfile="${MICROWHOISD_PIDFILE}"

depend() {
	need net
}
